# files/system/usr/lib/systemd/user/firefox-userjs-update.timer
[Unit]
Description=Daily update for Arkenfox user.js

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
