# files/system/usr/lib/systemd/user/firefox-theme-update.timer
[Unit]
Description=Daily update for Firefox GNOME Theme

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
